# supervisord-headless.conf
[unix_http_server]
file=/tmp/supervisor.sock

[supervisord]
nodaemon=true
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid
childlogdir=/tmp

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:chromium]
command=/app/start-chromium-headless.sh
autostart=true
autorestart=true
startsecs=3
startretries=5
stdout_logfile=/tmp/chromium-stdout.log
stderr_logfile=/tmp/chromium-stderr.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB

# socat for Chrome DevTools Protocol port forwarding
# Chromium's --remote-debugging-address=0.0.0.0 is disabled for security reasons,
# so we use socat to forward external connections to Chromium's localhost-only port.
# See: https://issues.chromium.org/issues/40261787
[program:socat]
command=socat TCP-LISTEN:9222,reuseaddr,fork TCP:localhost:9223
autostart=true
autorestart=true
startsecs=1
stdout_logfile=/tmp/socat-stdout.log
stderr_logfile=/tmp/socat-stderr.log
stdout_logfile_maxbytes=1MB
stderr_logfile_maxbytes=1MB
